# app/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
childlogdir=/var/log/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:celerycam]
command=/usr/bin/python manage.py celerycam
directory=/home/vagrant/code/
stdout_logfile=/var/log/celerycam.log
stderr_logfile=/var/log/celerycam.log
startsecs=5
stopwaitsecs=60
killasgroup=true
stopasgroup=true
priority=997

[program:celery]
command=/usr/bin/python manage.py celeryd -E -B --loglevel=INFO --verbosity=2
directory=/home/vagrant/code/
stdout_logfile=/var/log/celeryworker.log
stderr_logfile=/var/log/celeryworker.log
startsecs=5
stopwaitsecs=60
killasgroup=true
stopasgroup=true
priority=996

[program:celery-flower]
command=celery flower --broker_api=http://127.0.0.1:15672/api
directory=/home/vagrant/
stdout_logfile=/var/log/flower.log
stderr_logfile=/var/log/flower.log
startsecs=5
stopwaitsecs=60
killasgroup=true
stopasgroup=true
priority=998

[program:django-server]
command=/usr/bin/python manage.py runserver 0.0.0.0:8000
directory=/home/vagrant/code/
stdout_logfile=/var/log/djangoserver.log
stderr_logfile=/var/log/djangoserver.log
stopwaitsecs=2
killasgroup=true
stopasgroup=true
priority=999

[group:celeryapp]
programs=celerycam,celery,celery-flower,django-server
